# crontab/crontab_10
*/10 * * * * /usr/bin/php /var/www/html/tcc/classes/notificar.php

